# attn/attn_handler.service.in
[Unit]
Description=The Attention Handler service
Wants=obmc-mapper.target
After=obmc-mapper.target
After=mapper-wait@-xyz-openbmc_project-logging.service
After=fsi-scan@0.service
After=op-cfam-reset.service
After=cfam_override@0.service
After=openpower-update-bios-attr-table.service
Conflicts=op-stop-instructions@0.service
Conflicts=obmc-host-timeout@0.target

[Service]
@ENABLE_PHAL_TRUE@Environment="PDBG_DTB=/var/lib/phosphor-software-manager/hostfw/running/DEVTREE"
ExecStart=/usr/bin/openpower-hw-diags --daemon
Restart=on-failure

[Install]
WantedBy=obmc-host-startmin@0.target
